# resources/services/backup-cloud.path
[Unit]
Description=Nextcloud AIO backup lockfile watcher
After=network.service

[Path]
PathExists=/srv/nc-bkp/borg/lock.exclusive

[Install]
WantedBy=multi-user.target