# resources/services/backup-cloud.path
[Unit]
Description=Nextcloud AIO backup lockfile watcher
After=network.service

[Path]
PathChanged=/srv/nc-bkp/aio-lockfile

[Install]
WantedBy=multi-user.target